# cdh-ubuntu/cdh/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:cloudera-scm-agent]
command=/bin/bash -c "sudo service cloudera-scm-agent start"